# build/0.CORE/fdd/fs/etc/inittab
::sysinit:/etc/init
::askfirst:/bin/toybox
::restart:/bin/toybox reboot
::ctrlaltdel:/bin/toybox reboot
::shutdown:/bin/toybox umount -a -r
::ssh:/bin/dbclient
::profile:/etc/profile