# build/0.CORE/fdd/fs/etc/inittab
::sysinit:/etc/init.d/rc
::askfirst:/bin/toybox
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::ssh:/bin/dbclient
::sh:/bin/toybox sh
::bash:/bin/toybox sh